# raspberry_pi_zero2w_usb_camera_offline/chro-usb-camera.service
[Unit]
Description=CH-RO USB MJPEG Camera Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
SupplementaryGroups=video
WorkingDirectory=__INSTALL_DIR__
ExecStart=/usr/bin/python3 __INSTALL_DIR__/usb_mjpeg_server.py --host 0.0.0.0 --port 8080 --device /dev/video0 --width 640 --height 480 --fps 10 --jpeg-quality 70 --backend auto
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

